# ms-user/msuser.service
[Unit]
Description=service to handle users data

[Service]
Type=oneshot
ExecStart=/home/ubuntu/gosystemd-service/msuser
